# jibri-vm/systemd/jibri-autoscaler.timer
[Unit]
Description=Run Jibri Autoscaler every 20 seconds

[Timer]
OnBootSec=30
OnUnitActiveSec=20s
AccuracySec=5s
Unit=jibri-autoscaler.service

[Install]
WantedBy=timers.target
